# sup/hhvm.conf
[program:hhvm]
command = service hhvm start
user = nginx
autostart = true