# sup/hhvm.conf
[program:hhvm]
command = service hhvm start
user = root
autostart = true
