# sup/nginx.conf
[program:nginx]
command = service nginx start
user = nginx
autostart = true