# sup/php-fpm.conf
[program:php5-fpm]
command = service php5-fpm start
user = root
autostart = true
